# installer/wlcb-mixer.service
# =============================================================================
# WLCB-Mixer — systemd unit (atomic-release aware)
# =============================================================================
#
# This unit runs WLCB-Mixer from the *current* symlink:
#   /opt/wlcb-mixer/current -> /opt/wlcb-mixer/releases/<release-id>
#
# WHY THIS MATTERS
# ----------------
# Atomic release switching prevents "mixed version" bugs where the UI and server
# are from different commits. An upgrade becomes:
#   1) build new release in a new folder
#   2) flip symlink
#   3) restart service
#
# The service always runs a coherent set of files from one release directory.
# =============================================================================

[Unit]
Description=WLCB-Mixer Control Server
After=network.target

[Service]
User=wlcb
Group=wlcb

WorkingDirectory=/opt/wlcb-mixer/current/server
ExecStart=/usr/bin/node dist/index.js

Restart=always
RestartSec=5

Environment=NODE_ENV=production
EnvironmentFile=-/etc/wlcb-mixer/config.env

[Install]
WantedBy=multi-user.target
